# etc/supervisor/conf.d/sshd.conf
[program:test]
command=bash -c '/usr/sbin/sshd -D -e'
autostart=true
autorestart=true
stderr_logfile=/var/log/myapi.err.log
stdout_logfile=/var/log/myapi.out.log
directory=/root/
